# roles/scific/templates/crontab
1 1 * * * {{ docker_user_name }} docker start certbot renew
2 1 * * * {{ docker_user_name }} docker exec nginx -s reload
